# daemon/systemd/lap-registry-qm-init.service
[Unit]
Description=LightAP Service Registry QM Initializer
Documentation=file:///usr/share/doc/lap/com/SERVICE_DISCOVERY_ARCHITECTURE.md
Requires=lap-registry-qm.socket
After=lap-registry-qm.socket

[Service]
Type=exec
ExecStart=/home/ddk/1_workspace/2_middleware/LightAP/build/modules/Com/lap-registry-init --type=qm --socket=/run/lap/registry_qm.sock

# 标准输入/输出
StandardOutput=journal
StandardError=journal

# 资源限制
MemoryMax=50M
CPUQuota=20%
TasksMax=20

# 运行时目录
RuntimeDirectory=lap
RuntimeDirectoryMode=0755

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/lap

# 超时控制
TimeoutStartSec=10s
TimeoutStopSec=5s

# 重启策略
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
